# dot_config/autostart/keyring_unlock.desktop
[Desktop Entry]
Type=Application
Exec=.local/bin/keyring_unlock.sh
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Keyring Unlock
